# scripts/com.battlify.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.battlify.helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/battlify-helper</string>
        <string>daemon</string>
    </array>

    <!-- Start at boot and keep alive (re-launch if it crashes). -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Logs for diagnostics. -->
    <key>StandardErrorPath</key>
    <string>/var/log/battlify-helper.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/battlify-helper.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
